# nfs-antivirus-0.1/src/bin/detect_usb/detect_usb.service
[Unit]
Description=antivirus start detect usb

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/detect_usb start
ExecStop=/usr/libexec/detect_usb stop

[Install]
WantedBy=multi-user.target
